# fortune@.service
[Unit]
Description=Fortune quotes

[Service]
ExecStart=/usr/games/fortune
StandardOutput=socket
DynamicUser=yes
ProtectHome=true
PrivateUsers=true
